# openrc scripts/backlight
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="backlight daemon"
description="Turn Thinkpad backlight on at level 1"
command=/home/sleepy/scripts/startup/backlight.sh
command_args="${backlight_args}"
