# openrc scripts/scheduler
#!/sbin/openrc-run
# Copyright 2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

name="scheduler daemon"
description="Set default scheduler to mq-deadline at startup"
command=/home/sleepy/scripts/startup/scheduler.sh
command_args="${scheduler_args}"
